# overlay/system/etc/init/goodix-disable.conf
description "Goodix Key Disabler"

start on started lightdm

respawn

exec /usr/bin/goodix-key-disable
